# systemd/axentx-pain-validator-daemon.service
[Unit]
Description=axentx pain validator — cross-source confirmation gate between research and BD
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-pain-validator-daemon.py
Restart=always
RestartSec=15
MemoryMax=128M
CPUQuota=15%

[Install]
WantedBy=multi-user.target
