# systemd/axentx-secret-watchdog.timer
[Unit]
Description=Daily axentx secret-rotation watchdog
Requires=axentx-secret-watchdog.service

[Timer]
# Daily at 09:00 UTC (= 16:00 ICT). Random delay so simultaneous timer
# storms at the top of the hour don't collide on tiny e2-micro CPU.
OnCalendar=*-*-* 09:00:00 UTC
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
