# systemd/axentx-incident-responder-daemon.service
[Unit]
Description=axentx incident responder — auto-rerun GH Actions + Render redeploys
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-incident-responder-daemon.py
Restart=always
RestartSec=15
MemoryMax=128M
CPUQuota=20%

[Install]
WantedBy=multi-user.target
