# systemd/axentx-secret-watchdog.service
[Unit]
Description=axentx secret-rotation watchdog (oneshot — pair with .timer)
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=REPO_ROOT=/opt/surrogate-1-harvest
ExecStart=/bin/bash /opt/surrogate-1-harvest/bin/axentx-secret-watchdog.sh
StandardOutput=journal
StandardError=journal
