# systemd/surrogate-self-heal-daemon.service
[Unit]
Description=Self-heal daemon (continuous polling, restart dead services)
After=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/oci-self-heal-daemon.py
Restart=always
RestartSec=15
MemoryMax=64M
